# configs/systemd/btrfs-autodefrag.service
[Unit]
Description=Monthly BTRFS defragmentation
Documentation=man:btrfs

[Service]
ExecStart=btrfs filesystem defragment -r /
Type=simple

[Install]
WantedBy=multi-user.target
